# roles/extractor/scripts/systemd/exocortex-extractor.service
[Unit]
Description=Exocortex Extractor inbox-check
After=network.target

[Service]
Type=oneshot
ExecStart={{WORKSPACE_DIR}}/FMT-exocortex-template/roles/extractor/scripts/extractor.sh inbox-check
StandardOutput=append:{{HOME_DIR}}/logs/extractor/systemd-inbox-check.log
StandardError=append:{{HOME_DIR}}/logs/extractor/systemd-inbox-check-error.log
Environment=HOME={{HOME_DIR}}
Environment=PATH=/usr/local/bin:/usr/bin:/bin:/usr/local/sbin
