# roles/extractor/scripts/launchd/com.extractor.inbox-check.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.extractor.{{NAMESPACE}}.inbox-check</string>

    <key>ProgramArguments</key>
    <array>
        <string>{{ROOT_DIR}}/roles/extractor/scripts/extractor.sh</string>
        <string>--root-dir</string>
        <string>{{ROOT_DIR}}</string>
        <string>--workspace-dir</string>
        <string>{{WORKSPACE_DIR}}</string>
        <string>--agent-ai-path</string>
        <string>{{CLAUDE_PATH}}</string>
        <string>inbox-check</string>
    </array>

    <key>StartInterval</key>
    <integer>10800</integer>

    <key>StandardOutPath</key>
    <string>{{WORKSPACE_DIR}}/logs/extractor/launchd-inbox-check.log</string>

    <key>StandardErrorPath</key>
    <string>{{WORKSPACE_DIR}}/logs/extractor/launchd-inbox-check-error.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/opt/homebrew/bin</string>
    </dict>

    <key>RunAtLoad</key>
    <false/>
</dict>
</plist>
