# roles/synchronizer/scripts/systemd/exocortex-scheduler.service
[Unit]
Description=Exocortex Scheduler (central dispatcher)
After=network.target

[Service]
Type=oneshot
ExecStart={{WORKSPACE_DIR}}/{{EXOCORTEX_REPO}}/roles/synchronizer/scripts/scheduler.sh dispatch
StandardOutput=append:{{HOME_DIR}}/.local/state/logs/synchronizer/systemd-scheduler.log
StandardError=append:{{HOME_DIR}}/.local/state/logs/synchronizer/systemd-scheduler-error.log
Environment=HOME={{HOME_DIR}}
Environment=PATH=/usr/local/bin:/usr/bin:/bin:/usr/local/sbin
