# roles/synchronizer/scripts/systemd/exocortex-scheduler.timer
[Unit]
Description=Exocortex Scheduler timer — {{NAMESPACE}} (10 dispatches/day)

[Timer]
OnCalendar=*-*-* 00:00:00
OnCalendar=*-*-* 03:00:00
OnCalendar=*-*-* {{TIMEZONE_HOUR}}:00:00
OnCalendar=*-*-* 06:00:00
OnCalendar=*-*-* 09:00:00
OnCalendar=*-*-* 12:00:00
OnCalendar=*-*-* 15:00:00
OnCalendar=*-*-* 18:00:00
OnCalendar=*-*-* 21:00:00
OnCalendar=*-*-* 23:00:00
Persistent=true
Unit=exocortex-scheduler-{{NAMESPACE}}.service

[Install]
WantedBy=timers.target
