# roles/strategist/scripts/systemd/exocortex-strategist-morning.service
[Unit]
Description=Exocortex Strategist morning (day-plan / session-prep)
After=network.target

[Service]
Type=oneshot
ExecStart={{WORKSPACE_DIR}}/FMT-exocortex-template/roles/strategist/scripts/strategist.sh morning
StandardOutput=append:{{HOME_DIR}}/logs/strategist/systemd-morning.log
StandardError=append:{{HOME_DIR}}/logs/strategist/systemd-morning-error.log
Environment=HOME={{HOME_DIR}}
Environment=PATH=/usr/local/bin:/usr/bin:/bin:/usr/local/sbin
